# scripts/youtube-url-updater.task.xml
<?xml version="1.0" encoding="UTF-8"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Fires at login and hourly, but run_pipeline.py only does work once ~24h has passed since it last actually ran (gate measured from real run time, not a clock slot). So it runs about once a day, anchored to when you actually use the PC. Discovers new videos and fetches up to MAX_PER_RUN transcripts. The once-a-day cadence lets YouTube's per-IP rate-limit window reset between runs.</Description>
    <URI>\YouTube URL Updater</URI>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <UserId>YOUR_USER</UserId>
    </LogonTrigger>
    <TimeTrigger>
      <StartBoundary>2026-06-01T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <Repetition>
        <Interval>PT1H</Interval>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
    </TimeTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>YOUR_USER</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <Enabled>true</Enabled>
    <StartWhenAvailable>true</StartWhenAvailable>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <ExecutionTimeLimit>PT15M</ExecutionTimeLimit>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Windows\System32\wsl.exe</Command>
      <Arguments>-d Ubuntu -e bash -lc "cd /home/YOUR_USER/path/to/claudex-setup &amp;&amp; .venv-linux/bin/python scripts/run_pipeline.py"</Arguments>
    </Exec>
  </Actions>
</Task>
